# private_dot_config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=%E/neomutt/sync-mail.sh
