# private_dot_config/private_systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -Vac %E/isync/mbsyncrc
ExecStartPost=/usr/bin/notmuch new
