# private_dot_config/systemd/user/gpg-agent-browser.socket.d/override.conf
[Socket]
ListenStream=
ListenStream=%t/gnupg/d.uf4qnedu6tfpxtmo6gyoytsi/S.gpg-agent.browser